# deploy/solinteg-telemetry.timer
[Unit]
Description=Solinteg telemetry capture (10-min plan refresh)

[Timer]
# Every 10 minutes. Cheap (a localhost render: upstream fetches are cached, the DP runs in
# milliseconds), and caps plan staleness at ~10 min in the slow-drift band the dispatch
# loop's SoC-drift replan trigger can't see (drift below DISPATCH_REPLAN_DRIFT_KWH would
# otherwise ride until the next render). Telemetry growth is contained by logOptimizerRun's
# dedup + 3-decimal rounding (lib/telemetry/dispatch.ts) — an unchanged plan is not
# re-logged, so quiet ticks are nearly free.
OnCalendar=*:00/10
# Plus a few minutes past STOCKHOLM midnight (explicit timezone — the NUC's system clock
# may not be Stockholm): the render is what produces the new day's price snapshot +
# optimizer plan, and without this entry the dispatch loop sits planless (safely in auto,
# but silent and unoptimized) for up to an hour into every new day — the "midnight lull"
# visible as reduced armed_fraction in oracle_daily and priced into its regret.
OnCalendar=*-*-* 00:03:00 Europe/Stockholm
# Twice shortly after the ~13:00 Stockholm day-ahead release: the hourly :00 render races the
# release, and lib/prices.ts caches a today-only fetch with ~20-min rechecks, so without these
# the first unattended plan to see tomorrow's prices lands ~14:00 — afternoon plans decide
# evening hold-vs-sell on a truncated horizon for up to an hour. 13:22 clears the 13:00
# render's recheck window for a punctual release; 13:42 catches a late one. (Both become
# redundant if replanning ever gets an invokable trigger instead of riding renders.)
OnCalendar=*-*-* 13:22:00 Europe/Stockholm
OnCalendar=*-*-* 13:42:00 Europe/Stockholm
# If the NUC was powered off at a scheduled time, run once on next boot to catch up.
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
