# deploy/solinteg-boot-notify.service
[Unit]
Description=Notify on boot (resilience: proves auto-recovery actually happened)
After=network-online.target

[Service]
Type=oneshot
User=solinteg
WorkingDirectory=/opt/solinteg/app
EnvironmentFile=/opt/solinteg/solinteg.env
ExecStart=/opt/solinteg/app/.venv/bin/python /opt/solinteg/app/scripts/services/boot_notify.py
# Sandbox: no writes needed at all, unlike the other timer units.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
